# config/sign-nvidia.service
[Unit]
Description=Auto Sign NVIDIA Kernel Modules on Boot
Documentation=man:sign-nvidia-modules(8)
After=dracut-pre-build.service
After=network-online.target
Wants=network-online.target
ConditionVirtualization=!container

[Service]
Type=oneshot
ExecStart=/usr/local/bin/sign-nvidia-modules.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nvidia-signing
RemainAfterExit=yes
ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=false
PrivateTmp=yes

# Security hardening
ReadWritePaths=/usr/lib/modules /etc/pki/akmods/certs /var/lib/nvidia-signing /var/log/nvidia-signing /boot

[Install]
WantedBy=multi-user.target
